# auto-upload-ip/auto-upload-ip.service
[Unit]
Descripttion=Get IP Address
Wants=getip.timer

[Service]
ExecStart=/bin/getip.sh

[Install]
WantedBy=muti-user.target
